# systemd/bird-notify.service
[Unit]
Description=Bird detection → ntfy notification watcher
After=birdnet-go.service
Requires=birdnet-go.service

[Service]
Type=simple
User=sam
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 /home/sam/bird-listener/scripts/notify.py \
    --log /home/sam/bird-listener/logs/detections.log \
    --ntfy-url https://ntfy.sh/sam_b_bird_alerts
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
